# deploy/openclaw-load-campaign.timer
[Unit]
Description=Run OpenClaw Load Campaign at 5:30,6:30,7:30,8:30 AM US Eastern

[Timer]
OnCalendar=*-*-* 05,06,07,08:30:00
Timezone=America/New_York
Persistent=true
Unit=openclaw-load-campaign.service

[Install]
WantedBy=timers.target
